# templates/supervisor/thelounge.conf
[program:thelounge]
command = /usr/local/bin/docker-compose up
user = {{ ansible_user }}
directory = {{ thelounge_home }}
redirect_stderr = true
